# overlay/etc/cron.d/docker-clean
0 * * * * root docker images | awk '{print $3}' | xargs -n1 docker rmi