# overlay/etc/profile.d/etcd.sh
export ETCDCTL_CA_FILE=/etc/ssl/etcd/ca.pem
export ETCDCTL_CACERT=/etc/ssl/etcd/ca.pem
export ETCDCTL_ENDPOINT=https://localhost:2379
export ETCDCTL_ENDPOINTS=https://localhost:2379
export ETCDCTL_USER="$USER:$PASS"
