# overlay/etc/profile.d/etcd.sh
export ETCDCTL_CA_FILE=/etc/ssl/etcd/ca.pem
export ETCDCTL_CACERT=/etc/ssl/etcd/ca.pem